# files/usr/etc/profile.d/aliases.sh
# Aliases
alias ls="ls -hN --color=auto --group-directories-first"
alias ll="ls -lh --color=auto --group-directories-first"
alias l="ls -lah --group-directories-first"
alias la="ls -ah --group-directories-first"
alias cl="clear"
alias lc="ls -lZh --color=auto --group-directories-first"
alias lca="ls -lZah --color=auto --group-directories-first"
export LS_COLORS='rs=0:di=01;34:ln=01;36:mh=00:pi=40;33'
alias wg="wget"
